# whisper-dictate.service
[Unit]
Description=Whisper-Dictate daemon
After=graphical-session.target sound.target

[Service]
Type=simple
WorkingDirectory=%h/Code/whisper-dictate
ExecStart=%h/Code/whisper-dictate/.venv/bin/python %h/Code/whisper-dictate/dictate.py --compute-type int8 --cpu-threads 6 --beam-size 1 --no-vad-filter --no-condition-on-previous-text --no-speech-threshold 0.6
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
