# whisper-dictate-hotkey.service
[Unit]
Description=Whisper-Dictate global hotkey listener
After=graphical-session.target whisper-dictate.service
Wants=whisper-dictate.service

[Service]
Type=simple
WorkingDirectory=%h/Code/whisper-dictate
ExecStart=/usr/bin/python3 %h/Code/whisper-dictate/kglobal_hotkey.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
