# packaging/kdictate-systemd.service
[Unit]
Description=KDictate core daemon
After=graphical-session.target sound.target

[Service]
Type=simple
WorkingDirectory=@@REPO_DIR@@
ExecStart=@@REPO_DIR@@/.venv/bin/kdictate-daemon@@BACKEND_FLAGS@@
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
